# squid/modules/deploy_squid/install_bash/squid.service
[Unit]
Description=Squid Web Proxy.
[Service]
Type=simple
PIDFile=/var/run/squid.pid
ExecStart=/usr/sbin/squid
[Install]
WantedBy=multi-user.target